# UnifiOS/etc/cron.d/dns4me_cron
0 0 * * * root /data/custom/dns4me/dns4me.sh